# flask.conf
[program:api_freelance]
command = gunicorn api-freelance:app -b 0.0.0.0:5000
directory = /var/www/pos-sample
user = root
autostart=true
autorestart=true

